# airplay2-ledfx/services/shairport/run
#!/command/with-contenv bashio
bashio::log.info "SP: Begin shairport-sync..."

IFACE=$(bashio::config 'network_interface')
AIRPLAY_NAME=$(bashio::config 'airplay_name')

export PULSE_SERVER="unix:/run/pulse/native"
export DBUS_SYSTEM_BUS_ADDRESS="unix:path=/run/dbus-private/system_bus_socket"

until [ -S /run/pulse/native ]; do
	bashio::log.info "SP: Waiting for PulseAudio..."
	sleep 1
done

until [ -S /run/dbus-private/system_bus_socket ]; do
	bashio::log.info "SP: Waiting for D-Bus..."
	sleep 1
done

# Patch interface into shairport config at runtime
sed -i "s/interface = \"enu1\"/interface = \"${IFACE}\"/" /etc/shairport-sync.conf
sed -i "s/name = \"Home Assistant\"/name = \"${AIRPLAY_NAME}\"/" /etc/shairport-sync.conf

bashio::log.info "SP: Starting shairport-sync on interface ${IFACE} as '${AIRPLAY_NAME}'..."
if bashio::config.true 'verbose'; then
	bashio::log.info "SP: Version: $(shairport-sync --version 2>&1)" 
	bashio::log.info "SP: PULSE_SERVER: ${PULSE_SERVER}"
	bashio::log.info "SP: Pulse socket exists: $([ -S /run/pulse/native ] && echo YES || echo NO)"
	bashio::log.info "SP: pactl test:"
	pactl -s unix:/run/pulse/native info 2>&1 | head -5 | while read line; do bashio::log.info "SP: $line"; done

	bashio::log.info "SP: shairport-sync backends:"
	shairport-sync -h 2>&1 | grep -i "backend\|audio\|\-o" | while read line; do bashio::log.info "SP: $line"; done

#	shairport-sync --configfile=/etc/shairport-sync.conf -o pulseaudio -vvv --logOutputLevel \
	shairport-sync --configfile=/etc/shairport-sync.conf -o pulseaudio -v --logOutputLevel \
		2>&1 | while read line; do bashio::log.info "SP: $line"; done
else
	exec shairport-sync --configfile=/etc/shairport-sync.conf -o pulseaudio
fi
