# airplay-ledfx/services/avahi/run
#!/command/with-contenv bashio
bashio::log.info "AV: Begin Avahi..."

IFACE=$(bashio::config 'network_interface')
export DBUS_SYSTEM_BUS_ADDRESS="unix:path=/run/dbus-private/system_bus_socket"

sed -i "s/allow-interfaces=.*/allow-interfaces=${IFACE}/" /etc/avahi/avahi-daemon.conf

until [ -S /run/dbus-private/system_bus_socket ]; do
	bashio::log.info "AV: Waiting for D-Bus..."
	sleep 1
done

bashio::log.info "AV: Starting Avahi..."
exec avahi-daemon --no-drop-root --no-chroot
